# debian/bc-vigil.service
[Unit]
Description=BitCrafts Vigil - file integrity scanner web UI
Documentation=file:///usr/share/doc/bc-vigil/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bc-vigil
Group=bc-vigil
WorkingDirectory=/var/lib/bc-vigil
EnvironmentFile=/etc/bc-vigil/bc-vigil.env
ExecStart=/opt/venvs/bc-vigil/bin/bc-vigil
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/bc-vigil
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

AmbientCapabilities=CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_DAC_READ_SEARCH

[Install]
WantedBy=multi-user.target
